# npm/services.d/tailscale/finish
# File: services.d/tailscale/finish
#!/command/with-contenv bash

echo "[Tailscale] Service exited with status $1"